# srcs/Grafana/srcs/supervisord.conf
[supervisord]
nodaemon=true
user=root
[program:grafana]
command=./grafana/bin/grafana-server -homepath grafana
autorestart=true
[program:telegraf]
command=./telegraf/usr/bin/telegraf
autorestart=true